# pi-gen-config/stage-carpi/03-carpi-services/files/carpi.service
[Unit]
Description=CarPi OBD2 Dashboard
Documentation=file:///opt/carpi/README.md

# Start after X11 is running, Bluetooth is ready, and rfcomm is bound
After=network.target bluetooth.target carpi-rfcomm.service carpi-x11.service hostapd.service dnsmasq.service
Wants=bluetooth.target hostapd.service dnsmasq.service
Requires=carpi-rfcomm.service carpi-x11.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/carpi/carpi

ExecStart=/usr/bin/python3 /opt/carpi/carpi/main.py

# Restart on crash — critical for an always-on embedded device
Restart=always
RestartSec=5

# X11 display for pywebview/GTK/WebKit rendering
Environment=DISPLAY=:0

# Unbuffered Python output so logs appear in journalctl in real time
Environment=PYTHONUNBUFFERED=1

# Groups needed: video (framebuffer), bluetooth (rfcomm), dialout (/dev/rfcomm0)
SupplementaryGroups=video bluetooth dialout render input

# Kill cleanly on systemctl stop
TimeoutStopSec=10
KillSignal=SIGTERM

SyslogIdentifier=carpi

[Install]
WantedBy=multi-user.target
