# pi-gen-config/stage-carpi/03-carpi-services/files/carpi.service
[Unit]
Description=CarPi OBD2 Dashboard
Documentation=file:///opt/carpi/README.md

# Start after Bluetooth is ready and the rfcomm bind has run
After=network.target bluetooth.target carpi-rfcomm.service hostapd.service dnsmasq.service
Wants=bluetooth.target hostapd.service dnsmasq.service
Requires=carpi-rfcomm.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/carpi/carpi

ExecStart=/usr/bin/python3 /opt/carpi/carpi/main.py

# Restart on crash — critical for an always-on embedded device
Restart=always
RestartSec=5

# Kivy display: use DRM/KMS backend directly, no X11 required
Environment=KIVY_BCM_DISPMANX_ID=2
Environment=DISPLAY=

# Unbuffered Python output so logs appear in journalctl in real time
Environment=PYTHONUNBUFFERED=1

# Groups needed: video (framebuffer), bluetooth (rfcomm), dialout (/dev/rfcomm0)
SupplementaryGroups=video bluetooth dialout

# Kill cleanly on systemctl stop
TimeoutStopSec=10
KillSignal=SIGTERM

SyslogIdentifier=carpi

[Install]
WantedBy=multi-user.target
