# pi-gen-config/stage-carpi/03-carpi-services/files/carpi-rfcomm.service
[Unit]
Description=CarPi rfcomm bind (OBD2 Bluetooth serial port)
After=bluetooth.target
Wants=bluetooth.target
# Run before the main CarPi app
Before=carpi.service

[Service]
Type=oneshot
# Read OBD_MAC from config — extract it with Python so we don't duplicate it
ExecStartPre=/bin/bash -c 'OBD_MAC=$(python3 -c "import sys; sys.path.insert(0, \"/opt/carpi/carpi\"); import config; print(config.OBD_MAC)"); echo "rfcomm binding MAC: $OBD_MAC"'
ExecStart=/bin/bash -c '\
    OBD_MAC=$(python3 -c "import sys; sys.path.insert(0, \"/opt/carpi/carpi\"); import config; print(config.OBD_MAC)"); \
    CHANNEL=$(python3 -c "import sys; sys.path.insert(0, \"/opt/carpi/carpi\"); import config; print(config.OBD_BT_CHANNEL)"); \
    rfcomm bind /dev/rfcomm0 "$OBD_MAC" "$CHANNEL" || true'

# Remain active so carpi.service can depend on it
RemainAfterExit=yes

# If rfcomm fails, carpi.service still starts — it handles reconnection internally
StandardOutput=journal
StandardError=journal
SyslogIdentifier=carpi-rfcomm

[Install]
WantedBy=multi-user.target
