# pi-gen-airplay/stage-airplay/03-signalkit-services/files/signalkit-x11.service
[Unit]
Description=SignalKit X11 Display Server
Documentation=file:///opt/signalkit/README.md
After=systemd-user-sessions.service plymouth-quit.service
Wants=plymouth-quit.service

# Only start if there's a display device
ConditionPathExists=/dev/fb0

[Service]
Type=simple
# X must run as root to open /dev/tty1 and the GPU device nodes.
# It drops privileges internally. Clients (pywebview) connect over the
# X socket as the pi user — no root access is granted to them.

# Start a minimal X server on vt7 with no access control
# vt7 avoids conflicts with the kernel console on vt1
# -nocursor hides the mouse cursor (touchscreen/no-mouse kiosk)
# -ac disables host-based access control so the pi user can connect
ExecStart=/usr/bin/X :0 -nolisten tcp -nocursor -ac vt7 -dpms
# Wait for X to accept connections, then disable screen blanking and DPMS.
# ExecStartPost runs in parallel with ExecStart, so we must poll for readiness.
ExecStartPost=/bin/bash -c '\
  for i in $(seq 1 20); do \
    xdpyinfo -display :0 >/dev/null 2>&1 && break; \
    sleep 0.5; \
  done; \
  DISPLAY=:0 xsetroot -solid "#111111" 2>/dev/null; \
  DISPLAY=:0 xset s off 2>/dev/null; \
  DISPLAY=:0 xset -dpms 2>/dev/null; \
  DISPLAY=:0 xset s noblank 2>/dev/null; \
  echo 0 > /sys/module/kernel/parameters/consoleblank 2>/dev/null; \
  setterm --blank 0 --powerdown 0 > /dev/tty1 2>/dev/null; \
  setterm --blank 0 --powerdown 0 > /dev/tty7 2>/dev/null; \
  while true; do \
    sleep 60; \
    DISPLAY=:0 xset s off 2>/dev/null; \
    DISPLAY=:0 xset -dpms 2>/dev/null; \
    DISPLAY=:0 xset s noblank 2>/dev/null; \
  done &'
Restart=always
RestartSec=3

# X manages its own VT — no TTY directives needed.
# (TTYVHangup sends SIGHUP which kills X immediately after it attaches.)
StandardOutput=null
StandardError=journal

SyslogIdentifier=signalkit-x11

[Install]
WantedBy=multi-user.target
